# extra/cf-dns-auto-updater@.service
[Unit]
Description=Update a set of domain name records on Cloudflare
Requires=network.target

[Service]
Type=oneshot
User=%i
Group=%i
WorkingDirectory=/home/%i/.local/share/cf-dns-auto-updater
ExecStart=/usr/bin/cf-dns-auto-updater

[Install]
WantedBy=multi-user.target
